# poka.conf
description "Gunicorn application server running poka webpage"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid root
setgid www-data

env PATH=/root/.python_virtual_envs/sunrin-commute/bin
chdir /root/vhosts/poka
exec gunicorn --workers 3 --bind unix:poka.sock -m 007 app2:app
